# linux/preview/examples/mssql-mlservices/overlay/tmp/supervisord.conf
[supervisord]
childlogdir=/tmp
logfile=/tmp/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=trace
nodaemon=true

[program:sqlservr]
command=/opt/mssql/bin/sqlservr
user=root
stopasgroup=true
retries=3

[program:launchpadd]
command=/opt/mssql/bin/launchpadd -usens=false
stopasgroup=true
retries=3
